# docker/app/supervisor.conf
[program:auth-api]
user=root
directory=/root/auth-api/auth_api/
command=uwsgi --http :8000 --module auth_api.wsgi
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0